# src/templates/x11/image/etc/init.d/osjs-gentoo
#!/sbin/runscript
#
# OS.js Daemon
#
# Copyright (c) 2013 Anders Evenrud <andersevenrud@gmail.com>
# Released under the 2-clause BSD license.
#

pidfile="/var/run/osjs"
cwd="/usr/local/bin"
command="osjs-daemon"
user="osjs"

depend()
{
  need localmount
}

start()
{
  ebegin "Starting OS.js"
  start-stop-daemon --start --exec $command \
    --pidfile $pidfile -d $cwd -m -b --name OSjs
  eend $?
}

stop()
{
  ebegin "Stopping OS.js"
  start-stop-daemon --stop --exec $command \
    --pidfile $pidfile
  eend $?
}

